# conf/sresar.cron
* * * * * root /usr/local/lib/os_health/ssar/healing.sh > /dev/null 2>&1
